# data/openrc/gnome-session-wayland
#!/bin/sh
#
# gnome-session-wayland   GNOME Session (Wayland)
#
# chkconfig: 3 90 10
# description: GNOME Session (Wayland)
# pidfile: /run/user/<uid>/gnome-session-wayland.pid

### BEGIN INIT INFO
# Provides: gnome-session-wayland
# Required-Start: dbus gnome-shell-wayland gnome-settings-daemon-wayland gnome-session-dbus gnome-session-init gnome-session-monitor
# Required-Stop: gnome-session-monitor gnome-session-init gnome-session-dbus
# Default-Start: 3
# Default-Stop: 0 1 6
# Short-Description: GNOME Session (Wayland)
### END INIT INFO

GNOME_SESSION_EXEC="/usr/bin/gnome-session"
PIDFILE="${XDG_RUNTIME_DIR}/gnome-session-wayland.pid"

start() {
    echo "Starting GNOME Session (Wayland)..."

    if [ -f "$PIDFILE" ] && kill -0 "$(cat "$PIDFILE")" 2>/dev/null; then
        echo "Already running."
        return 0
    fi

    # Start dependencies first (simplified)
    for svc in gnome-shell-wayland gnome-settings-daemon-wayland \
               gnome-session-dbus gnome-session-init gnome-session-monitor; do
        if [ -x "/etc/xdg/gnome/$svc" ]; then
            /etc/xdg/gnome/$svc start
        fi
    done

    # Start GNOME session
    "${GNOME_SESSION_EXEC}" --session=wayland &
    echo $! > "$PIDFILE"
    echo "GNOME Session (Wayland) started with PID $(cat "$PIDFILE")"
}

stop() {
    echo "Stopping GNOME Session (Wayland)..."
    if [ -f "$PIDFILE" ]; then
        PID=$(cat "$PIDFILE")
        kill "$PID" 2>/dev/null
        rm -f "$PIDFILE"
    fi

    # Stop dependencies (reverse order)
    for svc in gnome-session-monitor gnome-session-init gnome-session-dbus \
               gnome-settings-daemon-wayland gnome-shell-wayland dbus; do
        if [ -x "/etc/rc.d/$svc" ]; then
            /etc/rc.d/$svc stop
        fi
    done
}

status() {
    if [ -f "$PIDFILE" ] && kill -0 "$(cat "$PIDFILE")" 2>/dev/null; then
        echo "GNOME Session (Wayland) running with PID $(cat "$PIDFILE")"
        return 0
    else
        echo "GNOME Session (Wayland) is not running."
        return 3
    fi
}

restart() {
    stop
    start
}

case "$1" in
    start)   start ;;
    stop)    stop ;;
    restart) restart ;;
    status)  status ;;
    *)       echo "Usage: $0 {start|stop|restart|status}" ; exit 1 ;;
esac

exit 0
