# data/openrc/gsd-xsettings
#!/bin/sh
#
# Generic SysVinit script for GNOME session daemons
#
# chkconfig: 3 90 10
# description: ${RC_SVCNAME} daemon
# pidfile: /run/user/<uid>/${RC_SVCNAME}.pid

### BEGIN INIT INFO
# Provides: ${RC_SVCNAME}
# Required-Start: dbus
# Required-Stop:
# Default-Start: 3
# Default-Stop: 0 1 6
# Short-Description: ${RC_SVCNAME} daemon
### END INIT INFO

DAEMON="/usr/libexec/${RC_SVCNAME}"
PIDFILE="${XDG_RUNTIME_DIR}/${RC_SVCNAME}.pid"
ARGS=""

# Optional environment overrides
[ -n "$DISPLAY" ] || DISPLAY=:0
export DISPLAY

[ -n "$DBUS_SESSION_BUS_ADDRESS" ] && export DBUS_SESSION_BUS_ADDRESS

start() {
    echo "Starting ${RC_SVCNAME}..."
    if [ -f "$PIDFILE" ] && kill -0 "$(cat "$PIDFILE")" 2>/dev/null; then
        echo "${RC_SVCNAME} already running with PID $(cat "$PIDFILE")"
        return 0
    fi

    # Start daemon in background
    "$DAEMON" $ARGS &
    echo $! > "$PIDFILE"
    echo "${RC_SVCNAME} started with PID $(cat "$PIDFILE")"
}

stop() {
    echo "Stopping ${RC_SVCNAME}..."
    if [ -f "$PIDFILE" ]; then
        kill "$(cat "$PIDFILE")" 2>/dev/null
        rm -f "$PIDFILE"
        echo "${RC_SVCNAME} stopped"
    else
        echo "${RC_SVCNAME} not running"
    fi
}

status() {
    if [ -f "$PIDFILE" ] && kill -0 "$(cat "$PIDFILE")" 2>/dev/null; then
        echo "${RC_SVCNAME} running with PID $(cat "$PIDFILE")"
        return 0
    else
        echo "${RC_SVCNAME} is not running"
        return 3
    fi
}

restart() {
    stop
    start
}

case "$1" in
    start)   start ;;
    stop)    stop ;;
    restart) restart ;;
    status)  status ;;
    *)       echo "Usage: $0 {start|stop|restart|status}" ; exit 1 ;;
esac

exit 0
